# build/supervisor/conf/kafka-connect.conf
[program:splunk-kafka-connect]
command=/opt/splunk-kafka-connect/start.sh

stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
stdout_logfile=/var/log/supervisor/%(program_name)s_stdout.log
stdout_logfile_maxbytes=0 ; disable rotation by supervisord use logrotate
stderr_logfile_maxbytes=0 ; disable rotation by supervisord use logrotate
stdout_logfile_backups=0  ; disable rotation by supervisord use logrotate
stderr_logfile_backups=0  ; disable rotation by supervisord use logrotate

autostart=true
autorestart=true

startretries=0
retries =0
