# build/cron/crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=

0 0 * * * root /usr/sbin/logrotate /etc/logrotate.conf > /dev/null
